# packaging/app-resource-report.cron
*/5 * * * * root /usr/sbin/resource2db >/dev/null 2>&1
